# pp.conf
[program:pp]
command=python3 /home/curator2/pp/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/pp/test.err.log
stdout_logfile=/home/curator2/pp/test.out.log